# relay/relay.service
# docushark-relay systemd unit template.
#
# Install:
#   1. Copy the relay binary to /usr/local/bin/relay (cargo build
#      --release && sudo install -m 0755 target/release/relay
#      /usr/local/bin/).
#   2. Create the relay user + data directory:
#        sudo useradd --system --home /var/lib/docushark-relay \
#             --shell /usr/sbin/nologin docushark-relay
#        sudo install -d -m 0750 -o docushark-relay -g docushark-relay \
#             /var/lib/docushark-relay
#   3. Roll a fresh config:
#        sudo -u docushark-relay /usr/local/bin/relay init \
#             --config /var/lib/docushark-relay/relay.toml
#   4. Copy this file to /etc/systemd/system/docushark-relay.service,
#      then:
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now docushark-relay

[Unit]
Description=DocuShark Relay (collaboration / MCP / auth server)
Documentation=https://github.com/JPE-Net-Technologies/docushark
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=docushark-relay
Group=docushark-relay
WorkingDirectory=/var/lib/docushark-relay
ExecStart=/usr/local/bin/relay serve --config /var/lib/docushark-relay/relay.toml --data-dir /var/lib/docushark-relay
Restart=on-failure
RestartSec=5s

# `relay serve` blocks on Ctrl-C / SIGINT and stops both listeners
# cleanly. systemd sends SIGTERM by default; the binary treats it the
# same way thanks to tokio::signal::ctrl_c() catching both on Unix.
KillSignal=SIGTERM
TimeoutStopSec=15s

# --- Hardening (modest; the relay only needs network + its data dir) ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/docushark-relay
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
